# systemd/swarm-discovery.service
[Unit]
Description=Swarm Discovery Service

[Service]
ExecStart=/usr/local/bin/swarm-discovery.sh