# systemd/swarm-manager-sync.service
[Unit]
Description=Sync manager config

[Service]
Type=oneshot
ExecStart=/usr/local/bin/swarm-manager-sync.sh